# cron/learnintouch.cron
00 01 * * * /usr/local/learnintouch/engine/cron/backup_all.sh

