# learnintouch/learnintouch.cron
00 01 * * * /usr/local/learnintouch/engine/cron/backup-all-websites.sh >> /usr/local/learnintouch/logs/cron.log 2>&1

