# learnintouch/learnintouch.cron
* * * * * echo 'Hello Supercronic' >> /usr/local/learnintouch/logs/cron.log 2>&1
00 01 * * * /usr/local/learnintouch/engine/cron/backup-all-websites.sh

